# ollama-proxy.service
[Unit]
Description=Auth proxy in front of local Ollama
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
User=YOUR_LINUX_USERNAME
WorkingDirectory=/home/YOUR_LINUX_USERNAME/tg-business-ai-bot
EnvironmentFile=/home/YOUR_LINUX_USERNAME/tg-business-ai-bot/.env
ExecStart=/home/YOUR_LINUX_USERNAME/tg-business-ai-bot/venv/bin/python ollama_proxy.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
